# setup/esb-scraper.service
[Unit]
Description=ESB Energy Scraper
After=network.target

[Service]
Type=oneshot
User=esb
WorkingDirectory=/opt/smart-meter-hdf-tracker-ha
ExecStart=/opt/smart-meter-hdf-tracker-ha/venv/bin/python -m scraper.download
StandardOutput=journal
StandardError=journal
