# setup/esb-api.service
[Unit]
Description=ESB Energy Tracker API
After=network.target

[Service]
Type=simple
User=esb
WorkingDirectory=/opt/esb-energy
ExecStart=/opt/esb-energy/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
